# deploy/polymarket-arb-shadow.service
[Unit]
Description=polymarket-arb shadow service

[Service]
WorkingDirectory=/opt/polymarket-arb
ExecStart=/usr/bin/env PYTHONPATH=polymarket_arb python3 -m polymarket_arb.app.supervisor_cli --root /var/lib/polymarket-arb --limit 1
Restart=always

[Install]
WantedBy=multi-user.target
